# scripts/ingest/ocr-worker.service
[Unit]
Description=Yggdrasil OCR Worker — procesa inbox/ocr/raw/ automáticamente
After=network.target

[Service]
Type=simple
User=%i
EnvironmentFile=%h/yggdrasil-dew/.env
Environment=YGGDRASIL_ROOT=%h/yggdrasil-dew
ExecStart=/bin/bash -c 'while true; do for f in "$YGGDRASIL_ROOT"/inbox/ocr/raw/*; do [ -f "$f" ] && bash "$YGGDRASIL_ROOT"/scripts/ingest/ocr-ingest.sh "$f" && mv "$f" "$YGGDRASIL_ROOT"/inbox/ocr/processed/; done; sleep 30; done'
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
